# etc/inittab
::sysinit:/etc/init/rcS
tty1::askfirst:-/bin/login
tty2::askfirst:-/bin/login
::ctrlaltdel:/etc/init/rc0
