# etc/rc.local
#!/bin/sh

/local/rescue/busybox syslogd

echo 1 > /sys/devices/system/cpu/microcode/reload

irqbalance &

sleep 10 && rmdir /srv &

ifconfig wlp2s0 down
macchanger -r wlp2s0
ifconfig wlp2s0 up

