# etc/rc.local
#!/bin/sh

irqbalance &

sleep 10 && rmdir /srv &
